# env/roles/energycoin/templates/energycoin-supervisor.conf
[program:energycoind]
user={{ ansible_env.USER }}
command={{ ansible_env.HOME }}/energycoin/src/energycoind -conf={{ ansible_env.HOME }}/.EnergyCoin/EnergyCoin.conf -datadir={{ ansible_env.HOME }}/.EnergyCoin

stderr_logfile={{ ansible_env.HOME }}/logs/energycoind.log
stdout_logfile={{ ansible_env.HOME }}/logs/energycoind.log